# com.vitalcity.evidence-radar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Weekly City Evidence Radar run. Install:
     cp com.vitalcity.evidence-radar.plist ~/Library/LaunchAgents/
     launchctl load ~/Library/LaunchAgents/com.vitalcity.evidence-radar.plist
     Runs Mondays 08:00. Edit RADAR_SITE_URL to the live URL after deploy. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vitalcity.evidence-radar</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/joshgreenman/Experiments/city-evidence-radar/run.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>RADAR_SITE_URL</key>
    <string>https://vitalcity-nyc.github.io/city-evidence-radar/</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/joshgreenman/Experiments/city-evidence-radar</string>
  <key>StandardOutPath</key>
  <string>/Users/joshgreenman/Experiments/city-evidence-radar/logs/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/joshgreenman/Experiments/city-evidence-radar/logs/launchd.err</string>
</dict>
</plist>
